# app/supervisord.conf
[supervisord]
nodaemon=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

[program:cron]
command=cron -f
autostart=true
autorestart=true
stdout_logfile=/var/log/cron/cron_stdout.log
stderr_logfile=/var/log/cron/cron_stderr.log
user=root

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx/nginx_stdout.log
stderr_logfile=/var/log/nginx/nginx_stderr.log
user=root

[program:php-fpm]
command=/usr/sbin/php-fpm8.2 -F
autostart=true
autorestart=true
stdout_logfile=/var/log/php8.2-fpm.log
stderr_logfile=/var/log/php8.2-fpm.log
user=root

#[program: kafka-config]
#command=python3 /kafka-producer-consumer/config.py
#stdout_logfile=/dev/stdout
#stderr_logfile=/dev/stderr
#autostart=true
#autorestart=false
#startretries=0


#[program: send-consumer]
#command=python3 /kafka-producer-consumer/send-consumer.py
#stdout_logfile=/dev/stdout
#stderr_logfile=/dev/stderr



#[program: receive-consumer]
#command:python3 /kafka-producer-consumer/receive-consumer.py
#stdout_logfile=/dev/stdout
#stderr_logfile=/dev/stderr
